# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio

echo "Starting Graphite Exporter..."
exec python3 -m graphite_api
